# setup/librarian.conf
[program:lazylibrarian-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/librarian.sh

